# rootfs.old/etc/services.d/readsb/finish
#!/usr/bin/execlineb -S0
# shellcheck shell=bash

s6-svc -d /var/run/s6/services/readsb
