# rootfs/etc/services.d/tar1090/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

sleep 5

srcdir=/var/readsb/readsb
service=tar1090

/usr/local/share/tar1090/tar1090.sh \
  "/run/$service" \
  "$srcdir" \
  2>&1 | mawk -W interactive '{print "[tar1090] " $0}'
